# systemd/StockWeekly.timer
[Unit]
Description=Call weekly stock service.

[Timer]
OnCalendar=weekly
Persistent=true
Unit=StockWeekly.service

[Install]
WantedBy=timers.target
